# nodo-monitoreo/nodo-monitoreo/systemd/nodo-monitoreo.service
# Servicio systemd: arranca el firmware automaticamente al encender la Pi
# y lo reinicia si se cae.
#
# Instalacion:
#   sudo cp systemd/nodo-monitoreo.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now nodo-monitoreo
#
# Ver el log en vivo:
#   journalctl -u nodo-monitoreo -f

[Unit]
Description=Nodo de Monitoreo Ambiental VitiScience (Grupo P8)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=plant123
# Necesario para acceder a los pines GPIO sin ser root
SupplementaryGroups=gpio

WorkingDirectory=/home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo/src
ExecStart=/home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo/venv/bin/python \
          /home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo/src/main.py \
          /home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo/config.ini
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
