# nodo-monitoreo/nodo-monitoreo/systemd/nodo-monitoreo.service
# Servicio systemd: arranca el firmware automaticamente al encender la Pi
# y lo reinicia si se cae.
#
# Instalacion:
#   sudo cp systemd/nodo-monitoreo.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now nodo-monitoreo
#
# Ver el log en vivo:
#   journalctl -u nodo-monitoreo -f

[Unit]
Description=Nodo de Monitoreo Ambiental VitiScience (Grupo P8)
After=bluetooth.service ModemManager.service
Wants=bluetooth.service ModemManager.service

[Service]
Type=simple

WorkingDirectory=/home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo
ExecStart=/home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo/venv/bin/python /home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo/src/main.py /home/plant123/nodo-5G-G8/nodo-monitoreo/nodo-monitoreo/config.ini

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
